# pcp/config/rc.local
#!/usr/bin/sh
. $PCP_DIR/etc/pcp.env

case $1 in
	'start')
		# Install ZFS PMDA if not already registered
		if ! grep -q "^zfs" $PCP_PMCDCONF_PATH 2>/dev/null; then
			cd /var/lib/pcp/pmdas/zfs && echo "both" | ./Install >/dev/null 2>&1
		fi
		;;
	'stop')
		;;
esac
